# udev/99-usb-serial.rules
KERNEL=="ttyACM*", ATTRS{idVendor}=="TODO", SYMLINK+="emaks/can_usb", MODE="0666", RUN+="/usr/local/bin/slcan.sh"